# data/data/agent/systemd/units/node-zero.service
[Unit]
Description=Identify node zero to run OpenShift Assisted Installation Service on
Wants=network-online.target load-config-iso@.service
After=network-online.target load-config-iso@.service

[Service]
Type=oneshot
RemainAfterExit=True
ExecStart=/usr/local/bin/set-node-zero.sh

[Install]
WantedBy=multi-user.target
